# supervisord.conf
[supervisord]
nodaemon=true

[program:ser2net]
command=/usr/sbin/ser2net -n

